# docker/cron/crontab
0 0 * * * docker exec portfolio-django /django/scripts/logrotate.sh > /dev/stdout
